# launchd/com.rishi.trading-agent.decision-run.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Rendered by scripts/generate_plist.py into ~/Library/LaunchAgents (SPEC §9).
     launchd provides a minimal environment: PATH is set explicitly, and
     `claude` + Max OAuth must be verified under launchd before Phase 4 trusts
     this schedule. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{LABEL}}</string>
    <key>ProgramArguments</key>
    <array>
        <string>uv</string>
        <string>run</string>
        <string>python</string>
        <string>-m</string>
        <string>etrade_agent.runner</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{WORKDIR}}</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{PATH}}</string>
        <key>AUTOMATED</key>
        <string>1</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>{{HOUR}}</integer>
        <key>Minute</key>
        <integer>{{MINUTE}}</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>{{WORKDIR}}/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{WORKDIR}}/logs/launchd-stderr.log</string>
</dict>
</plist>
